# dist/usr/lib/systemd/system/coffer.socket
[Unit]
Description=Credential encryption FIFO

[Socket]
ListenStream=/run/coffer.socket
Accept=yes
