# dist/usr/lib/systemd/system/coffer@.service
[Unit]
Description=Encrypt and store credential

[Service]
Type=oneshot
StandardInput=socket
StandardOutput=socket
StandardError=journal
StateDirectory=coffer
DynamicUser=yes
ImportCredential=coffer.recipients.*
ImportCredential=coffer.identities.*
ExecStart=/usr/bin/coffer.py
